# send-request.timer
[Unit]
Description=HTTP Request Timer
[Timer]
OnBootSec=30min
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target
